# ops/systemd/hammer-radar/hammer-autonomous-trigger-scheduler-dry-run.timer.template
[Unit]
Description=Run Hammer Radar autonomous trigger scheduler dry-run tick every 2 minutes
Documentation=file:/home/josue/workspace/kernel/ai-agent-orchestrator-main/docs/hammer_radar/live_readiness/R290_MANUAL_SYSTEMD_DRY_RUN_TIMER_ACTIVATION_CHECKLIST.md

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
RandomizedDelaySec=5s
Persistent=false
Unit=hammer-autonomous-trigger-scheduler-dry-run.service

[Install]
WantedBy=timers.target
